# Workarounds/Performance Tuning/ac.target
[Unit]
Description=Runnning on AC power
DefaultDependencies=no
Conflicts=battery.target