# Workarounds/Performance Tuning/battery.target
[Unit]
Description=Running on Battery Power
DefaultDependencies=no
Conflicts=ac.target